# DEBIAN/postinst
#!/bin/sh -e
echo service template postinst
systemctl enable service-template.service
systemctl start service-template.service
